# templates/echo.service.erb
[Unit]
Description=echo
After=network.target

[Service]
ExecStart=<%= @exec_start %>
User=<%= @system_user %>
Restart=always

[Install]
WantedBy=multi-user.target
